# recipes-custom/services/files/cerber-client.service
[Unit]
Description=Cerber client service

[Service]
ExecStart=/usr/bin/cerber-client
Restart=always

[Install]
WantedBy=multi-user.target